# deploy/systemd/secretary-bot.service
[Unit]
Description=Telegram Secretary Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/secretary-bot/app
ExecStart=/opt/secretary-bot/venv/bin/python /opt/secretary-bot/app/main.py --config /opt/secretary-bot/runtime/config.yaml
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/opt/secretary-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
